# scripts/immich-backup.service
[Unit]
Description=Immich backup

[Service]
Type=oneshot
ExecStart=/bin/sh -lc '/usr/bin/env uv run %h/.local/bin/backup.py immich_postgres immich_server immich_machine_learning -D immich -U postgres --backup-dir /mnt/media/share/backups/databases/'
StandardOutput=journal
StandardError=journal
