# scripts/immich-backup.timer
[Unit]
Description=Immich backup
Requires=immich-backup.service

[Timer]
OnCalendar=*-*-* 02:58:49
Persistent=true

[Install]
WantedBy=timers.target
